# deploy/suape-frontend-deploy.timer
[Unit]
Description=Suape frontend auto-deploy poller (GHCR pull, every 2 minutes)

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
Unit=suape-frontend-deploy.service

[Install]
WantedBy=timers.target
